# clean-dollchan (1).service
[Unit]
Description=Clean dollchan backups

[Service]
Type=oneshot
ExecStart=/bin/bash %h/sync/main/Documents/dollchan/clean.sh
