# clean-sessions.timer
[Unit]
Description=daily clean tab sessions files

[Timer]
OnCalendar=daily
RandomizedDelaySec=12h
Persistent=true

[Install]
WantedBy=timers.target
