# clean-dollchan (1).timer
[Unit]
Description=Clean dollchan backups weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
