# clean-backups (1).service
[Unit]
Description=Clean system/home backups
Conflicts=backup_all.service

[Service]
Type=oneshot
ExecStart=/bin/bash %h/sync/system-data/clean.sh
