# fisher.timer
[Unit]
Description=Update fisher and plugins monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
