# google-docs (1).timer
[Unit]
Description=Download google documents daily

[Timer]
OnCalendar=*-*-* 16:00:00
Persistent=true

[Install]
WantedBy=timers.target
