# dotfiles_pi/service.conf/rngd/rngd.service.d/rngd-opts.conf
[Service]
EnvironmentFile=/etc/rngd.confd