# reflector/reflector.service.d/home.conf
[Service]
ExecStart=
ExecStart=/usr/bin/reflector --age 6 -p https -f 4 --save /etc/pacman.d/mirrorlist
