# systemd-units/lgogdownloader-update-cache.service
[Unit]
Description=Update gog games cache

[Service]
Type=oneshot
ExecStart=/usr/bin/lgogdownloader --update-cache
