# systemd-units/pihole-up.timer
[Unit]
Description=Run pihole update every two days

[Timer]
OnCalendar=daily
Persistent=true
 
[Install]
WantedBy=timers.target

