# systemd-units/xset-kbd-rate.service
[Unit]
Description=Set keyboard rate through xset

[Service]
Type=oneshot
ExecStart=/usr/bin/xset r rate 200 30

[Install]
WantedBy=graphical.target
