# systemd-units/xsetdpms.service
[Unit]
Description=Set dpms timeout through xset

[Service]
ExecStart=/usr/bin/xset dpms 600
Type=oneshot

[Install]
WantedBy=graphical.target
