# units/.config/systemd/user/backup.timer
[Unit]
Description=Backup everything daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
